# deploy/supervisord/meetpy.ini
[program:meetpy]
directory=/opt/meetpy/meetpy
command=gunicorn -t 30 meetpy.wsgi -b 127.0.0.1:8000
autostart=true